# config/includes.chroot/etc/xdg/autostart/anydesk.desktop
[Desktop Entry]
Name=AnyDesk
Terminal=false
Exec=anydesk
Type=Application
#Icon=/path/icon.png
Categories=Utility;
