# config/includes.chroot/etc/xdg/autostart/compton.desktop
[Desktop Entry]
Name=Compton
Description=Compton for X11
Terminal=false
Exec=compton
Type=Application
#Icon=/path/icon.png
Categories=Utility;
